# phosh.service.d/20-wlroots-hwc-env.conf
[Unit]
After=lxc@android.service

[Service]
Environment=WLR_BACKENDS='hwcomposer,libinput'
