# phosh.service.d/20-wlroots-hwc-env.conf
[Unit]
After=lxc@android.service

[Service]
Environment=WLR_BACKENDS='hwcomposer,libinput'
Environment=EGL_PLATFORM='hwcomposer'
Environment=WLR_HWC_SKIP_VERSION_CHECK=1
