# phosh.service.d/30-force-user.conf
[Service]
# FIXME: since uid 1000 is hardcoded several times, we need to force user on hybris adaptations
User=32011
