# phosh.service.d/40-droid-wait.conf
[Service]
# Wait for hwcomposer to start up...
ExecStartPre=/usr/bin/waitforservice init.svc.vendor.hwcomposer-*
